# kalium-iso-overlay/etc/sv/nix-daemon/run
#!/bin/sh
exec 2>&1
if [ -e /nix/var/nix/profiles/default/etc/profile.d/nix-daemon.sh ]; then
  . /nix/var/nix/profiles/default/etc/profile.d/nix-daemon.sh
fi

exec $(which nix-daemon)
